# deploy/systemd/echochamber-dc.target
[Unit]
After=network.target

Wants=echochamber-dc.service

[Install]
WantedBy=multi-user.target
